# filebrowser/filebrowser.service
[Unit]

Description=filebrowser

    

[Service]

User=root

ExecStart=/usr/bin/filebrowser --port 996

Restart=on-abort

    

[Install]

WantedBy=multi-user.target
